# deploy/crontab
*/5 * * * * %(user)s %(manage)s task_name
